# images/centos/seyren/centos6/etc/supervisor/conf.d/seyren.conf
[program:seyren]
command=java -jar /opt/seyren.jar

[program:mongodb]
command=/usr/bin/mongod -f /etc/mongodb.conf run
